# rootfs/etc/services.d/nordvpn/finish
#!/bin/bash
set -x

rm -rf /run/nordvpn
ip link delete nordlynx 2>/dev/null