# rootfs/etc/services.d/nordvpn/run
#!/bin/bash
set -x

if [[ ! -d /run/nordvpn ]]; then
  mkdir -m 0770 /run/nordvpn
fi

exec s6-notifyoncheck -n 60 -d /usr/sbin/nordvpnd > /tmp/nordvpnd.log 2>&1
